# standalone-apps/app-lab-desktop/build/debian/arduino-app-lab/etc/xdg/autostart/ArduinoAppLab.desktop
[Desktop Entry]
Encoding=UTF-8
Type=Application
Name=Arduino App Lab
Comment=
Exec=$BINARY_PATH
OnlyShowIn=XFCE;
RunHook=0
StartupNotify=false
Terminal=false
MimeType=x-scheme-handler/arduino-app-lab;
Hidden=false
Categories=Development;Utility;